# docs/infrastructure/scripts/parser-worker.conf
; Supervisor config for Laravel parser queue workers
; Copy to: /etc/supervisor/conf.d/parser-worker.conf

[program:parser-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/online-parser.siteaacess.store/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/online-parser.siteaacess.store/storage/logs/worker.log
stopwaitsecs=3600
